# templates/apt-mirror.service
[Unit]
Description=Ubuntu Selective Offline Mirror Materialize Service
After=network-online.target
Wants=network-online.target
Conflicts=apt-mirror-full.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ubuntu-offline-mirror.sh materialize-selective
StandardOutput=journal
StandardError=journal
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=infinity
KillMode=control-group
SendSIGKILL=no

[Install]
WantedBy=multi-user.target
