# templates/apt-mirror.service
[Unit]
Description=APT Mirror Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/ubuntu-mirror/run-apt-mirror.sh
StandardOutput=append:/var/log/apt-mirror.log
StandardError=append:/var/log/apt-mirror.log
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
