# templates/apt-mirror.service
[Unit]
Description=Ubuntu Offline Mirror Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ubuntu-offline-mirror.sh sync
StandardOutput=journal
StandardError=journal
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=infinity
KillMode=control-group
SendSIGKILL=no

[Install]
WantedBy=multi-user.target
